# anime_stock.service
[Unit]
Description=Animetrics AI Dashboard
Documentation=https://github.com/yourusername/anime-stock
After=network.target mysql.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/anime-stock
Environment="PATH=/opt/anime-stock/venv/bin:/usr/local/bin:/usr/bin"
EnvironmentFile=/opt/anime-stock/.env

# Run Streamlit
ExecStart=/opt/anime-stock/venv/bin/streamlit run \
    src/anime_stock/dashboard/app.py \
    --server.port=8501 \
    --server.address=127.0.0.1 \
    --server.headless=true \
    --browser.gatherUsageStats=false

# Restart on failure
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=append:/opt/anime-stock/logs/streamlit.log
StandardError=append:/opt/anime-stock/logs/streamlit-error.log

[Install]
WantedBy=multi-user.target
